# data/crontab
55 3 * * * cd /data &&bash /data/sv.sh c >> /data/cron.log
55 20 */6 * * cd /data &&bash /data/sv.sh s >> /data/cron.log
55 18 * * * cd /data &&bash /data/switch_node.sh >> /data/cron.log
